# mailgateway-integration/supervisor_delivery_watcher.conf
[program:delivery-watcher]
command=php /var/www/mailgateway/artisan delivery:watch --verbose
directory=/var/www/mailgateway
autostart=true
autorestart=true
startsecs=10
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/delivery-watcher.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30
stopsignal=SIGTERM

; Environment variables - set these for your environment
environment=
    RELAY_MONGODB_HOST="YOUR_RELAY_SERVER_IP",
    RELAY_MONGODB_PORT="27017",
    RELAY_MONGODB_DATABASE="relay_logs",
    RELAY_MONGODB_USERNAME="mailgateway_reader",
    RELAY_MONGODB_PASSWORD="YOUR_PASSWORD"
